# .config/systemd/user/urxvtd.service
[Unit]
Wants=xorg.target
After=xorg.target

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=xinit.target
