# .config/systemd/user/xinit.target
[Unit]
Description=Xinitrc Stuff
After=xorg.target
