# .config/systemd/user/xinit.target
[Unit]
Description=Xinitrc Stuff
Requires=xorg.target
After=xorg.target
