# .config/systemd/user/xkbmap.service
[Unit]
Wants=xorg.target
After=xorg.target

[Service]
ExecStart=/usr/bin/setxkbmap -layout gb

[Install]
WantedBy=xmonad.target
