# .config/systemd/user/xscreensaver.service
[Unit]
Description=XScreensaver
After=xorg.target

[Service]
ExecStart=/usr/bin/xscreensaver -no-splash

[Install]
WantedBy=xinit.target
