# .config/systemd/user/xsetroot.service
[Unit]
#Wants=xorg.target

[Service]
ExecStart=/usr/bin/xsetroot -cursor_name left_ptr

[Install]
WantedBy=xmonad.target
